# ansible/roles/neurondesktop/templates/neurondesk-api.service.j2
[Unit]
Description=NeuronDesktop API
After=network.target postgresql.service neuronagent.service
Requires=postgresql.service

[Service]
Type=simple
User={{ neurondb_user }}
Group={{ neurondb_group }}
WorkingDirectory={{ neurondesktop_bin_path }}
Environment="DB_HOST={{ neurondesktop_db_host }}"
Environment="DB_PORT={{ neurondesktop_db_port }}"
Environment="DB_NAME={{ neurondesktop_db_name }}"
Environment="DB_USER={{ neurondesktop_db_user }}"
Environment="DB_PASSWORD={{ neurondesktop_db_password }}"
Environment="SERVER_PORT={{ neurondesktop_api_port }}"
Environment="NEURONAGENT_ENDPOINT={{ neurondesktop_agent_endpoint }}"
Environment="NEURONAGENT_API_KEY={{ neurondesktop_agent_api_key }}"
Environment="LOG_LEVEL={{ neurondesktop_log_level }}"
ExecStart={{ neurondesktop_bin_path }}/neurondesk-api
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neurondesk-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target

