# ansible/roles/neurondesktop/templates/neurondesk-frontend.service.j2
[Unit]
Description=NeuronDesktop Frontend
After=network.target neurondesk-api.service
Requires=neurondesk-api.service

[Service]
Type=simple
User={{ neurondb_user }}
Group={{ neurondb_group }}
WorkingDirectory={{ neurondesktop_frontend_path }}
Environment="NEXT_PUBLIC_API_URL=http://localhost:{{ neurondesktop_api_port }}/api/v1"
Environment="PORT={{ neurondesktop_frontend_port }}"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neurondesk-frontend

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target

