# packaging/deb/neuronagent/package/etc/systemd/system/neuronagent.service
[Unit]
Description=NeuronAgent - AI Agent Runtime System
Documentation=https://www.neurondb.ai/docs/neuronagent
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=neuronagent
Group=neuronagent
WorkingDirectory=/var/lib/neuronagent
ExecStart=/usr/bin/neuronagent
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neuronagent

# Environment variables (can be overridden in /etc/neuronagent/config.yaml)
Environment="CONFIG_PATH=/etc/neuronagent/config.yaml"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/neuronagent /var/log/neuronagent

[Install]
WantedBy=multi-user.target

