# ansible/roles/neuronagent/templates/neuronagent.service.j2
[Unit]
Description=NeuronAgent - Agent runtime and REST API
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User={{ neurondb_user }}
Group={{ neurondb_group }}
WorkingDirectory={{ neuronagent_bin_path }}
ExecStart={{ neuronagent_bin_path }}/agent-server -config {{ neuronagent_config_dir }}/config.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neuronagent

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target

