# packaging/deb/neuronagent/DEBIAN/postinst
#!/bin/bash
set -e

# Post-installation script for NeuronAgent

# Create neuronagent user and group if they don't exist
if ! getent group neuronagent > /dev/null 2>&1; then
    groupadd --system neuronagent
fi

if ! getent passwd neuronagent > /dev/null 2>&1; then
    useradd --system --gid neuronagent --home-dir /var/lib/neuronagent \
            --shell /usr/sbin/nologin --comment "NeuronAgent service user" neuronagent
fi

# Create directories
mkdir -p /etc/neuronagent
mkdir -p /var/lib/neuronagent
mkdir -p /var/log/neuronagent

# Set permissions
chown -R neuronagent:neuronagent /var/lib/neuronagent
chown -R neuronagent:neuronagent /var/log/neuronagent
chmod 755 /etc/neuronagent
chmod 755 /var/lib/neuronagent
chmod 755 /var/log/neuronagent

# If config doesn't exist, copy example
if [ ! -f /etc/neuronagent/config.yaml ]; then
    if [ -f /etc/neuronagent/config.yaml.example ]; then
        cp /etc/neuronagent/config.yaml.example /etc/neuronagent/config.yaml
        chown neuronagent:neuronagent /etc/neuronagent/config.yaml
        chmod 640 /etc/neuronagent/config.yaml
        echo "Created /etc/neuronagent/config.yaml from example"
        echo "Please edit /etc/neuronagent/config.yaml to configure your database connection"
    fi
fi

# Reload systemd and enable service (but don't start)
if [ -f /etc/systemd/system/neuronagent.service ]; then
    systemctl daemon-reload
    systemctl enable neuronagent.service
    echo "NeuronAgent service enabled. Start it with: systemctl start neuronagent"
fi

echo "NeuronAgent installed successfully."
echo "Configuration: /etc/neuronagent/config.yaml"
echo "Logs: /var/log/neuronagent/"

exit 0



