# packaging/pkg/neuronmcp/scripts/postinstall
#!/bin/bash

# Post-install script for NeuronMCP macOS package

# Create neuronmcp user and group if they don't exist
if ! dscl . -read /Groups/neuronmcp > /dev/null 2>&1; then
    dscl . -create /Groups/neuronmcp
    dscl . -create /Groups/neuronmcp gid 501
fi

if ! dscl . -read /Users/neuronmcp > /dev/null 2>&1; then
    dscl . -create /Users/neuronmcp
    dscl . -create /Users/neuronmcp UserShell /usr/bin/false
    dscl . -create /Users/neuronmcp RealName "NeuronMCP service user"
    dscl . -create /Users/neuronmcp UniqueID 501
    dscl . -create /Users/neuronmcp PrimaryGroupID 501
    dscl . -create /Users/neuronmcp NFSHomeDirectory /var/lib/neuronmcp
fi

# Create directories
mkdir -p /var/lib/neuronmcp
mkdir -p /var/log/neuronmcp

# Set permissions
chown -R neuronmcp:neuronmcp /var/lib/neuronmcp
chown -R neuronmcp:neuronmcp /var/log/neuronmcp
chmod 755 /var/lib/neuronmcp
chmod 755 /var/log/neuronmcp

# If config doesn't exist, copy example
if [ ! -f /etc/neuronmcp/mcp-config.json ]; then
    if [ -f /etc/neuronmcp/mcp-config.json.example ]; then
        cp /etc/neuronmcp/mcp-config.json.example /etc/neuronmcp/mcp-config.json
        chown neuronmcp:neuronmcp /etc/neuronmcp/mcp-config.json
        chmod 640 /etc/neuronmcp/mcp-config.json
        echo "Created /etc/neuronmcp/mcp-config.json from example"
        echo "Please edit /etc/neuronmcp/mcp-config.json to configure your database connection"
    fi
fi

# Set permissions on launchd plist
if [ -f /Library/LaunchDaemons/com.neurondb.neuronmcp.plist ]; then
    chown root:wheel /Library/LaunchDaemons/com.neurondb.neuronmcp.plist
    chmod 644 /Library/LaunchDaemons/com.neurondb.neuronmcp.plist
fi

echo "NeuronMCP installed successfully."
echo ""
echo "Configuration: /etc/neuronmcp/mcp-config.json"
echo "Logs: /var/log/neuronmcp/"
echo ""
echo "To use with MCP clients (e.g., Claude Desktop), configure them to use:"
echo "  /usr/local/bin/neurondb-mcp"
echo ""
echo "Optional: To run as a background service:"
echo "  sudo launchctl load /Library/LaunchDaemons/com.neurondb.neuronmcp.plist"
echo "  sudo launchctl start com.neurondb.neuronmcp"

