# packaging/deb/neurondesktop/neurondesktop.service
[Unit]
Description=NeuronDesktop - Unified Web Interface for NeuronDB Ecosystem
Documentation=https://www.neurondb.ai/docs/neurondesktop
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=neurondesktop
Group=neurondesktop
WorkingDirectory=/var/lib/neurondesktop
ExecStart=/usr/bin/neurondesktop-api
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neurondesktop

# Environment variables (can be overridden in /etc/neurondesktop/config.yaml)
Environment="CONFIG_PATH=/etc/neurondesktop/config.yaml"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/neurondesktop /var/log/neurondesktop /var/www/neurondesktop

[Install]
WantedBy=multi-user.target


