# ansible/roles/neuronmcp/templates/neuronmcp.service.j2
[Unit]
Description=NeuronMCP - MCP Server for NeuronDB
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User={{ neurondb_user }}
Group={{ neurondb_group }}
WorkingDirectory={{ neuronmcp_bin_path }}
Environment="NEURONDB_HOST={{ neuronmcp_db_host }}"
Environment="NEURONDB_PORT={{ neuronmcp_db_port }}"
Environment="NEURONDB_DATABASE={{ neuronmcp_db_name }}"
Environment="NEURONDB_USER={{ neuronmcp_db_user }}"
Environment="NEURONDB_PASSWORD={{ neuronmcp_db_password }}"
Environment="NEURONMCP_HTTP_ADDR={{ neuronmcp_http_addr }}"
Environment="NEURONDB_LOG_LEVEL={{ neuronmcp_log_level }}"
Environment="NEURONDB_LOG_FORMAT={{ neuronmcp_log_format }}"
Environment="NEURONDB_LOG_OUTPUT={{ neuronmcp_log_output }}"
ExecStart={{ neuronmcp_bin_path }}/neurondb-mcp
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neuronmcp

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target

